# deploy/hoster-kg/streamlit-sales-analytics.service
[Unit]
Description=Sales Analytics Streamlit App
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/opt/sales-analytics
Environment="PATH=/opt/sales-analytics/.venv/bin"
ExecStart=/opt/sales-analytics/.venv/bin/python -m streamlit run app.py --server.port 8501 --server.address 127.0.0.1 --server.headless true
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
